# deluge.conf
[program:deluged]
command = /usr/bin/deluged -d -c /config/deluge -L info -l /config/deluge/deluged.log

[program:deluge-web]
command = /usr/bin/deluge-web -c /config/deluge
